# provisioning/roles/grafana/files/grafana.service
[Unit]
Description=Node Exporter

[Service]
ExecStart=/opt/monitoring/grafana-2.5.0/run.sh

[Install]
WantedBy=default.target
